# launchd/com.sefk.agentsview.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.sefk.agentsview</string>

	<!-- Keep agentsview serving + watching for new Claude/Ollama sessions. -->
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/agentsview</string>
		<string>serve</string>
		<string>--no-browser</string>
		<string>--host</string>
		<string>0.0.0.0</string>
		<string>--port</string>
		<string>8088</string>
	</array>

	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>

	<!-- launchd starts with a minimal PATH; give agentsview a sane one. -->
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>

	<key>StandardOutPath</key>
	<string>/Users/sefk/.agentsview/serve.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/sefk/.agentsview/serve.log</string>
</dict>
</plist>
